# src/httpd.ini
[program:httpd]
command=/usr/sbin/httpd -DFOREGROUND
stopsignal=TERM
stopwaitsecs=10
