# src/sssd.ini
[program:sssd]
command=/usr/sbin/sssd -i --logger=files
redirect_stderr=true
